# dist/read-temperature.service
[Unit]
Description=Monitors Room Temperature
After=influxdb.service
Requires=influxdb.service

[Service]
User=grayshade
ExecStart=/home/grayshade/read-temperature-tokio/target/release/read-temperature-tokio
Restart=on-failure
RestartSec=1

[Install]
WantedBy=default.target
